# templates/unicorn@.service.erb
# original: http://bogomips.org/unicorn.git/tree/examples/unicorn@.service?h=v5.0.1
# ==> /etc/systemd/system/unicorn@.service <==
# Since SIGUSR2 upgrades do not work under systemd, this service file
# allows starting two simultaneous services during upgrade time
# (e.g. unicorn@1 unicorn@2) with the intention that they take
# turns running in-between upgrades.  This should allow upgrading
# without downtime.

[Unit]
Description = unicorn Rack application server %i
Wants = unicorn.socket
After = unicorn.socket

[Service]
User = <%= @user %>
<% if @group -%>
Group = <%= @group %>
<% end -%>
<% if @working_directory -%>
WorkingDirectory = <%= @working_directory %>
<% end -%>
<% @environment.each do |key, value| -%>
Environment = "<%= "#{key}=#{value}" %>"
<% end unless @environment.empty? -%>
# bundler users must use the "--keep-file-descriptors" switch, here:
# ExecStart = bundle exec --keep-file-descriptors unicorn -c ...
ExecStart = <%= @exec_start %>
ExecReload = /bin/kill -HUP $MAINPID
KillSignal = SIGQUIT
Sockets = unicorn.socket

# This is based on the Unicorn::Configurator#timeout directive,
# adding a few seconds for scheduling differences:
TimeoutStopSec = 62

[Install]
WantedBy = multi-user.target
